# lib/etude_for_ops/templates/common/cloud/aws/99_share/platform/ruby/chef/templates_default_active_job.service.erb
[Unit]
Description=ActiveJob

[Service]
Environment=RAILS_ENV=<%= @rails_env %>
WorkingDirectory=/srv/<%= @app %>/current
ExecStart=/home/<%= @user %>/.rbenv/shims/bundle exec sidekiq -C /srv/<%= @app %>/shared/config/sidekiq.yml
ExecStop=kill -s SIGTERM `cat /srv/<%= @app %>/shared/pids/sidekiq.pid`
Restart=always
Type=simple
User=<%= @user %>
Group=<%= @group %>

[Install]
WantedBy = multi-user.target
