# lib/etude_for_ops/templates/common/cloud/aws/99_share/platform/ruby/chef/templates_default_kuroko2_executor.erb
#!/bin/bash
# chkconfig: 345 99 20
# description: Kuroko2 executor
# processname: kuroko2_executor

# Include RedHat function library
. /etc/rc.d/init.d/functions

NAME=kuroko2_executor
USER=<%= @user %>
APP_NAME=<%= @app %>
APP_DIR=<%= @app_path %>/current
SHARE_DIR=<%= @app_path %>/shared
PID_DIR=<%= @app_path %>/shared/pids
LOG_DIR=<%= @app_path %>/shared/log
BUNDLE=/home/$USER/.rbenv/shims/bundle
RAILS_ENV=<%= @rails_env %>

sudo -u $USER -H sh -c "mkdir -p $PID_DIR"
sudo -u $USER -H sh -c "mkdir -p $LOG_DIR"

PIDFILE=$PID_DIR/command-executor.pid


start() {
cd $APP_DIR
daemon --pidfile=$PIDFILE --user=$USER RAILS_ENV=$RAILS_ENV $BUNDLE exec rails runner 'Kuroko2::Servers::CommandExecutor.new.run' &
}

stop() {
echo -n "Stopping Kuroko2 Executor"
kill -s SIGTERM `cat $PIDFILE`
}

restart() {
stop
start
}

case "$1" in
start)
start
;;
stop)
stop
;;
restart)
restart
;;
*)
N=/etc/init.d/$NAME
echo "Usage: $N {start|stop|restart}" >&2
exit 1
;;
esac
