# lib/etude_for_ops/templates/common/cloud/aws/99_share/platform/ruby/chef/templates_default_puma.erb
#!/bin/bash
# chkconfig: 345 98 20
# description: Application server
# processname: puma

# Include RedHat function library
. /etc/rc.d/init.d/functions

NAME=puma
USER=<%= @user %>
APP_NAME=<%= @app %>
APP_DIR=<%= @app_path %>/current
SHARE_DIR=<%= @app_path %>/shared
PID_DIR=<%= @app_path %>/shared/pids
SOCKET_DIR=<%= @app_path %>/shared/system/sockets
LOG_DIR=<%= @app_path %>/shared/log
BUNDLE=/home/$USER/.rbenv/shims/bundle
RAILS_ENV=<%= @rails_env %>
S3_ACCESS_KEY=<%= @s3_access_key %>
S3_SECRET_KEY=<%= @s3_secret_key %>
SECRET_KEY_BASE=<%= @secret_key_base %>

sudo -u $USER -H sh -c "mkdir -p $PID_DIR"
sudo -u $USER -H sh -c "mkdir -p $SOCKET_DIR"
sudo -u $USER -H sh -c "mkdir -p $LOG_DIR"

PIDFILE=$PID_DIR/puma.pid
ULOCK=/var/lock/subsys/$NAME

start() {
cd $APP_DIR
# Start puma
echo -n $"Starting $NAME: "
export S3_ACCESS_KEY=$S3_ACCESS_KEY
export S3_SECRET_KEY=$S3_SECRET_KEY
export SECRET_KEY_BASE=$SECRET_KEY_BASE
daemon --pidfile=$PIDFILE --user=$USER $BUNDLE exec puma -C $SHARE_DIR/config/puma.rb -d
puma=$?
[ $puma -eq 0 ] && touch $ULOCK
echo

return $puma
}

stop() {
echo -n "Stopping puma"
kill -s SIGTERM `cat $PIDFILE`
}

restart() {
stop
start
}

case "$1" in
start)
start
;;
stop)
stop
;;
restart)
restart
;;
*)
N=/etc/init.d/$NAME
echo "Usage: $N {start|stop|restart}" >&2
exit 1
;;
esac
