# lib/etude_for_ops/templates/common/cloud/aws/99_share/platform/ruby/chef/templates_default_puma.service.erb
[Unit]
Description=puma

[Service]
Environment=RAILS_ENV=development
WorkingDirectory=/srv/<%= @app %>/current
ExecStart=/home/<%= @user %>/.rbenv/shims/bundle exec puma -C /srv/<%= @app %>/shared/config/puma.rb
ExecStop=kill -s SIGTERM `cat /srv/<%= @app %>/shared/pids/puma.pid`
Restart = always
Type = simple
User=<%= @user %>
Group=<%= @group %>

[Install]
WantedBy = multi-user.target
