# lib/etude_for_ops/templates/common/cloud/aws/99_share/platform/ruby/docker/ship/application_service/app/etc/supervisord.conf
[supervisord]
nodaemon=true

[program:zabbix-agent]
command=/bin/bash -c "service zabbix-agent start"

[program:fluentd]
command=/bin/bash -c "service td-agent start"

[program:rails]
command=/usr/local/bundle/bin/rails server -b 0.0.0.0
