# lib/etude_for_ops/templates/onpremis/99_share/platform/ruby/chef/templates_default_puma.service.erb
[Unit]
Description=puma

[Service]
Environment=RAILS_ENV=development
WorkingDirectory=/var/app/current
ExecStart=/home/vagrant/.rbenv/shims/bundle exec puma -C /var/app/shared/puma.rb
ExecStop=kill -s SIGTERM `cat /var/app/shared/tmp/pids/puma.pid`
Restart = always
Type = simple

[Install]
WantedBy = multi-user.target