# lib/etude_for_ops/templates/common/99_share/platform/ruby/chef/templates_default_td-agent-ui.service.erb
[Unit]
Description=puma

[Service]
Environment=RAILS_ENV=development
ExecStart=/usr/sbin/td-agent-ui start
ExecStop=/usr/sbin/td-agent-ui stop
Restart = always
Type = simple
User=<%= @user %>
Group=<%= @group %>

[Install]
WantedBy = multi-user.target