# lib/etude_for_ops/templates/common/cloud/aws/99_share/platform/ruby/docker/ship/application_service/proxy-app/etc/supervisord.conf
[supervisord]
nodaemon=true

[program:zabbix-agent]
command=/bin/bash -c "service zabbix-agent start"

[program:fluentd]
command=/bin/bash -c "service td-agent start"

[program:nginx]
command=/bin/bash -c "service nginx start"
